# deploy/llm-router.service
[Unit]
Description=FastAPI llm-router service
After=network.target

[Service]
User=zhouli
Group=zhouli
WorkingDirectory=/opt/llm_router
Environment="PATH=/opt/llm_router/.venv/bin"
ExecStart=/opt/llm_router/.venv/bin/gunicorn main:app -k uvicorn.workers.UvicornWorker -w 2 -b 127.0.0.1:8000 --access-logfile - --error-logfile -
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
